# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

[program:disco]
command=/usr/local/bin/disco nodaemon
autostart=true
autorestart=true
