# build/base/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:disco]
command=/usr/local/bin/disco nodaemon

